# services/waterfall.service
[Unit]
Description=waterfall
After=network.target
Wants=waterfall.timer

[Service]
Type=oneshot
ExecStart=/home/pi/poolpi/waterfall
